# containers/nginx/crontab
# Run logrotate every minute for nginx logs
* * * * * /usr/sbin/logrotate /etc/logrotate.d/nginx 